# firmware/raspberry_pi_gateway/systemd/plant-health-gateway.service
[Unit]
Description=Plant Health Gateway - Raspberry Pi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/plant_health_gateway
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/plant_health_gateway/src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
